# legacy/vulnscout-mcp.service
[Unit]
Description=VulnScout MCP Server
After=network.target

[Service]
Type=simple
User=kali
WorkingDirectory=/home/kali/vulnscout
EnvironmentFile=/home/kali/vulnscout/.env
Environment="VULNSCOUT_DIR=/home/kali/vulnscout"
ExecStart=/home/kali/vulnscout/.venv/bin/python /home/kali/vulnscout/mcp_server.py --port 8000
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
